# examples/task_third/supervisord/conf.d/crontab.conf
[program:cron]
command=cron -f -L 8
autostart=true
autorestart=true